# airootfs/etc/skel/.config/autostart/owlsetup.desktop
[Desktop Entry]
Name=Owl Setup
Exec=/home/guest/.startscript
Icon=utilities-terminal
Type=Application
Categories=System;Utility;
NoDisplay=true
